# com.user.logi-wake-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.logi-wake-watcher</string>
    <key>Program</key>
    <string>BINARY_PATH</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/logi-wake-watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/logi-wake-watcher.err</string>
</dict>
</plist>
